# anbox-openrc/anbox-session.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Anbox Session Manager"
supervisor="supervise-daemon"
command="/usr/bin/anbox"
command_args="session-manager"
